# crontab
0 7 * * * /Users/jasonamaya/Desktop/REMEDI-med-dispenser/center/sh/gitpull.sh
